# temp/ofono/ofono-auto-enable.initd
#!/sbin/openrc-run

command="/usr/bin/ofono-auto-enable"

depend() {
   need dbus
}
